# daemoneye.service
[Unit]
Description=DaemonEye AI Tmux Daemon
Documentation=https://github.com/daemoneye/daemoneye
After=network.target

[Service]
Type=simple
# Assuming the daemoneye binary is installed via cargo or placed in ~/.cargo/bin
ExecStart=%h/.cargo/bin/daemoneye daemon
ExecStop=%h/.cargo/bin/daemoneye stop
Restart=on-failure
RestartSec=5

# Ensure tmux can be found in the PATH
Environment="PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

# If you use environment variables for API keys, you can specify an EnvironmentFile:
# EnvironmentFile=%h/.config/daemoneye/.env

[Install]
WantedBy=default.target
