# t1000.service
[Unit]
Description=T1000 AI Tmux Daemon
Documentation=https://github.com/t1000/t1000
After=network.target

[Service]
Type=simple
# Assuming the t1000 binary is installed via cargo or placed in ~/.cargo/bin
ExecStart=%h/.cargo/bin/t1000 daemon
ExecStop=%h/.cargo/bin/t1000 stop
Restart=on-failure
RestartSec=5

# Ensure tmux can be found in the PATH
Environment="PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin"

# If you use environment variables for API keys, you can specify an EnvironmentFile:
# EnvironmentFile=%h/.config/t1000/.env

[Install]
WantedBy=default.target
